# tvheadend/rootfs/etc/s6-overlay/s6-rc.d/init-epg/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: TVHeadend
# Starts O2 server
# ==============================================================================

exec /usr/local/bin/update-epg.sh
