# oneplayserver/rootfs/etc/s6-overlay/s6-rc.d/init-epg/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: O2 OnePlay server
# Starts epg init
# ==============================================================================

# wait for the OnePlay service to start
while ! nc -z 0.0.0.0 8082; do
    sleep 1
done

exec /usr/local/bin/update-epg.sh
